# supervisord.conf
[supervisord]
nodaemon=true

[program:daphne]
command=poetry run daphne -b 0.0.0.0 -p 8001 crawlsy.asgi:application
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/daphne.err.log
stdout_logfile=/var/log/daphne.out.log

[program:celery_beat]
command=poetry run celery -A crawlsy beat -l info
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/celery_beat.err.log
stdout_logfile=/var/log/celery_beat.out.log

[program:celery_worker]
command=poetry run celery -A crawlsy worker -l info
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/celery_worker.err.log
stdout_logfile=/var/log/celery_worker.out.log
